# infra/quadlets/api.container
[Unit]
Description=My API
After=postgres.service valkey.service
Wants=postgres.service valkey.service

[Container]
Image=ghcr.io/OWNER/api:prod
ContainerName=api
Pod=orbitcheck.pod
EnvironmentFile=/srv/myapp/api/env/api.env
Volume=/var/lib/containers/volumes/api-data:/app/data:Z,U
Label="io.containers.autoupdate=image"
NoNewPrivileges=yes
ReadOnly=true
Tmpfs=/tmp:size=64M

[Service]
Restart=always
RestartSec=3
# CORRECTED: The command is now wrapped in a shell to allow redirection
ExecStartPre=-/bin/sh -c '/usr/bin/podman login ghcr.io -u your_github_username --password-stdin < /run/secrets/ghcr_token'

[Install]
WantedBy=default.target