# infra/quadlets/minio.container
[Unit]
Description=MinIO
After=network-online.target
Wants=network-online.target

[Container]
Remove=false
Image=docker.io/minio/minio:latest
ContainerName=minio
Pod=orbitcheck.pod
EnvironmentFile=/srv/myapp/minio/env/minio.env
Volume=/var/lib/containers/volumes/minio-data:/data:Z,U
# Console on 9001, S3 on 9000; expose only via Caddy unless you want public S3 ports
# Access internally within the Pod via service name "minio:9000"
Exec=server /data --console-address :9001
Label="io.containers.autoupdate=image"
NoNewPrivileges=yes

[Service]
Restart=always
ExecStartPre=mkdir -p /var/lib/containers/volumes/minio-data

[Install]
WantedBy=default.target