# infra/quadlets/promtail.container
[Unit]
Description=Promtail
After=loki.service
Wants=loki.service

[Container]
Remove=false
Image=docker.io/grafana/promtail:2.9.8
ContainerName=promtail
Pod=orbitcheck.pod
Volume=/var/lib/containers/volumes/promtail-config:/etc/promtail:Z
# For rootless, scrape container logs from $XDG_RUNTIME_DIR/containers or JSON logs
# Mount the containers storage if needed:
# Volume=/home/podmanuser/.local/share/containers:/containers:ro
Label="io.containers.autoupdate=image"
NoNewPrivileges=yes

[Service]
Restart=always
ExecStartPre=mkdir -p /var/lib/containers/volumes/promtail-config

[Install]
WantedBy=default.target