# infra/quadlets/statping.container
[Unit]
Description=Statping-NG

[Container]
Image=docker.io/statping-ng/statping-ng:v0.90.80
ContainerName=statping
Pod=orbitcheck
Volume=/var/lib/containers/volumes/statping-data:/app:Z,U
EnvironmentFile=/srv/myapp/statping/env/statping.env
ContainerAutoUpdate=registry
NoNewPrivileges=yes

[Service]
Restart=always
[Install]
WantedBy=default.target