# infra/quadlets/valkey.container
[Unit]
Description=Valkey
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/valkey/valkey:7.2.11
ContainerName=valkey
Pod=orbitcheck.pod
Volume=/var/lib/containers/volumes/valkey-data:/data:Z,U
Label="io.containers.autoupdate=image"
NoNewPrivileges=yes
ReadOnly=false

[Service]
Restart=always

[Install]
WantedBy=default.target