# infra/quadlets/api/api.container
[Unit]
Description=My API
After=api-db.service valkey.service
Wants=api-db.service valkey.service

[Container]
LogDriver=journald
Image=ghcr.io/guillemso1er/api:prod
ContainerName=api
Pod=orbitcheck.pod

EnvironmentFile=/home/podmanuser/.config/api/api.env
EnvironmentFile=/home/podmanuser/.config/api/api.secrets.env

Label="io.containers.autoupdate=image"
NoNewPrivileges=yes
ReadOnly=true
Tmpfs=/tmp:size=64M

[Service]
Restart=always
RestartSec=3

# Ensure env dir exists
ExecStartPre=/usr/bin/install -d -m 700 /var/lib/containers/configs/api/env

# Fetch secrets just-in-time, atomically, with strict perms
# Uses the service token from a local file (not shown in process args)
ExecStartPre=/bin/sh -ceu ' \
    tok="$(cat /home/podmanuser/.secrets/infisical/api.token)"; \
    export INFISICAL_TOKEN="$tok"; \
    /bin/infisical export \
        --domain=https://infisical.bastiat.xyz \
        --path=/api \
        --env=prod \
        --format=dotenv \
        --silent \
        --include-imports \
    > /var/lib/containers/configs/api/.api.secrets.env.tmp; \
    chmod 600 /var/lib/containers/configs/api/.api.secrets.env.tmp; \
    mv /var/lib/containers/configs/api/.api.secrets.env.tmp /home/podmanuser/.config/api/api.secrets.env \
'


[Install]
WantedBy=default.target